# dwm/local/systemd/user/dwmbar.timer
[Unit]
Description=Update DWM bar regularly

[Timer]
OnCalendar=*:*:00
AccuracySec=500ms

[Install]
WantedBy=timers.target