# nvidia/system/lib/systemd/system/nvidia-disable.service
[Unit]
Description=Disables Nvidia GPU

[Service]
Type=oneshot
ExecStart=-/usr/sbin/nvidia-stop

[Install]
WantedBy=multi-user.target
